# app/Resources/fixtures/crontab
# vim: filetype=crontab

# Location of the log directory

# Location of the PHP Command Line Interface binary.

# Export tweetwallPro
* * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php export_tweetwall_pro 2>&1 | logify >>$LOGDIR/project_export_tweetwall_pro.log 2>&1

# IndexContent & images
* * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_admin index_content --site=project_admin >/dev/null 2>&1
* * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_lelab_admin index_content --site=project_lelab_admin >/dev/null 2>&1
* * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_images index_images >/dev/null 2>&1

# Programmated publication
* * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_admin unstack_publication_pending --action=create 2>&1 | logify >>$LOGDIR/project_unstack_publish_create.log 2>&1
* * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_admin unstack_publication_pending --action=delete 2>&1 | logify >>$LOGDIR/project_unstack_publish_delete.log 2>&1

# Chamane
25 6 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:import:chamane --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_import_chamane.log 2>&1
28 6 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:import:playlist-live --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_playlist_live.log 2>&1
0 10 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:import:musical-references --siteaccess=project -i 1day --env=prod 2>&1 | logify >>$LOGDIR/project_musical_references.log 2>&1
30 4 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:project:radioline --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_radioline.log 2>&1

# Believe + Dailymotion
* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:project:believe-video-upload --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_believe_upload.log 2>&1
*/2 * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_admin videoUpload --service=Dailymotion 2>&1 | logify >>$LOGDIR/project_dailymotion_upload.log 2>&1

# Google editor choice
*/15 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console project:flux:google-editor-choice --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_google_choice.log 2>&1

# Twitter
* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:project:twitter-export --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_twitter_export.log 2>&1

# Factiva & Cedrom
*/5 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:project:migration-export --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_export_factiva.log 2>&1
*/5 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:project:migration-export --siteaccess=project --type=cedrom --env=prod 2>&1 | logify >>$LOGDIR/project_export_cedrom.log 2>&1

# eZFlow
* * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php ezflow 2>&1 | logify >>$LOGDIR/ezflow.log 2>&1
* * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php ezflow-cleanup 2>&1 | logify >>$LOGDIR/ezflow-cleanup.log 2>&1

# Sitemaps (general & images)
0 23 * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_lelab generate_sitemap_general 2>&1 | logify >>$LOGDIR/lelab-generate_sitemap_general.log 2>&1
0 23 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:images-feed --env=prod 2>&1 | logify >>$LOGDIR/lelab-generate_sitemap_images.log 2>&1

# Sitemap (news)
*/10 * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_lelab generate_sitemap_news 2>&1 | logify >>$LOGDIR/lelab-generate_sitemap_news.log 2>&1

# Génération ez
*/10 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:tagsHeaderFooter --siteaccess=project_lelab --env=prod 2>&1 | logify >>$LOGDIR/lelab-generate_html_part.log 2>&1
#Replace by above*/10 * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_lelab generate_html_part 2>&1 | logify >>$LOGDIR/lelab-generate_html_part.log 2>&1

# Génération symfony
0 9,16 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:zone-file --env=prod --siteaccess=project_lelab left 2>&1 | logify >>$LOGDIR/project_lelab-left.log 2>&1
1 9,16 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:zone-file --env=prod --siteaccess=project_lelab left_article 2>&1 | logify >>$LOGDIR/project_lelab-left_article.log 2>&1
*/10 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:xml-articles --env=prod -m 50 2>&1 | logify >>$LOGDIR/project_lelab-xml-articles.log 2>&1
2 9,16 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:tops-ten --env=prod --siteaccess=project_lelab 2>&1 | logify >>$LOGDIR/project_lelab-tops.log 2>&1
*/10 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:feed-list --env=prod 2>&1 | logify >>$LOGDIR/project_lelab-feed-list.log 2>&1
* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:eztag-items-table --env=prod 2>&1 | logify >>$LOGDIR/project_lelab-eztags-items-table.log 2>&1
* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:laselectedtweets --env=prod 2>&1 | logify >>$LOGDIR/project_lelab-laselectedtweets.log 2>&1
43 5-22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console social:update:raisedItems --siteaccess=project_lelab --env=prod 2>&1 | logify >>$LOGDIR/project_lelab-raisedItems.log 2>&1
3 9,16 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console social:update:raisedItems --siteaccess=project_lelab --env=prod pending 2>&1 | logify >>$LOGDIR/project_lelab-raisedItems_pending.log 2>&1
52 5-22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console social:update:raisedItems --siteaccess=project --interval=3days --env=prod 2>&1 | logify >>$LOGDIR/project-raisedItems.log 2>&1

* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:headersListening --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_headers_listening.log 2>&1
*/5 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:highlights:header --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_highlights_header.log 2>&1
*/5 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:highlights:footer --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_highlights_footer.log 2>&1
5 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:most-commented --siteaccess=project --dateInterval=2days --env=prod 2>&1 | logify >>$LOGDIR/project_most_commented.log 2>&1
10 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:most-read --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_most_read.log 2>&1
*/5 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:rightcol --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_rightcol.log 2>&1
*/2 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:timeline --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_timeline.log 2>&1
* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:blocks --siteaccess=project top --env=prod 2>&1 | logify >>$LOGDIR/project_blocks_top.log 2>&1
*/5 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:blocks --siteaccess=project bottom --env=prod 2>&1 | logify >>$LOGDIR/project_blocks_bootom.log 2>&1
*/10 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:feed-list --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_feed_list.log 2>&1
*/30 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:most-shared --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_most_shared.log 2>&1
15 23 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console project:generate:top-archive --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_top_archives.log 2>&1
*/30 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:best-page-hub-tag --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_best_page_hub.log 2>&1
20 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:top-page-hub --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_top_page_hub.log 2>&1
*/5 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:direct-last-content --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_direct_last_content.log 2>&1
* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:send:contactform --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_send_contactform.log 2>&1
* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:newsSitemap --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_sitemap_news.log 2>&1
8 22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:sitemapPictures --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_sitemap_pictures.log 2>&1
9 22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:classicSitemap --siteaccess=project --type=speaker --env=prod 2>&1 | logify >>$LOGDIR/project_sitemap_speaker.log 2>&1
11 22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:classicSitemap --siteaccess=project --type=page_hub --env=prod 2>&1 | logify >>$LOGDIR/project_sitemap_pagehub.log 2>&1
12 22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:sitemapVideos --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_sitemap_video.log 2>&1
16 22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:indexSitemap --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_sitemap_index.log 2>&1
17 22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:articleByCategorySitemap --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_article_category.log 2>&1
37 22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:sitemapMediaCenter --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_sitemap_mediacenter.log 2>&1
* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:project:home-generation --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_home_generation.log 2>&1
36 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:project:active-broadcast --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_get_active_broadcast.log 2>&1
#* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:googleNewsRSS --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_gnews.log 2>&1
* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:send:soundArchivesUser --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_send_sound_archives.log 2>&1
35 22 * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_admin generate_dashboard_analytics --type=production 2>&1 | logify >>$LOGDIR/project-generate_dashboard_analytics.log 2>&1
0 23 * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_admin remove_media_files 2>&1

* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console club:home-generation --siteaccess=project_club --env=prod 2>&1 | logify >>$LOGDIR/project_club_home_generation.log 2>&1
* * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:home-generation --siteaccess=project_lelab --env=prod 2>&1 | logify >>$LOGDIR/project_lelab_home_generation.log 2>&1
37 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:rab-de-lab --siteaccess=project_lelab --env=prod 2>&1 | logify >>$LOGDIR/project_lelab_rab_de_lab.log 2>&1
38 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:more-home --siteaccess=project_lelab --env=prod 2>&1 | logify >>$LOGDIR/project_lelab_more_home.log 2>&1
*/10 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:news-feed --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_news_feed.log 2>&1
39 22 * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_admin generate_dashboard 2>&1 | logify >>$LOGDIR/project_dashboard.log 2>&1
* * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_lelab_admin clear_memcache_pending --site=project_lelab_admin 2>&1 | logify >>$LOGDIR/project_lelab_clear_memcache.log 2>&1
* * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_admin clear_memcache_pending --site=project_admin 2>&1 | logify >>$LOGDIR/project_clear_memcache.log 2>&1
52 22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:highlightsList --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_highlights_list.log 2>&1

38 22 * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_lelab_admin generate_dashboard_analytics --type=production 2>&1 | logify >>$LOGDIR/project_lelab-generate_dashboard_analytics.log 2>&1
43 22 * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_lelab_admin generate_dashboard 2>&1 | logify >>$LOGDIR/project_lelab_dashboard.log 2>&1
0 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console meteo:generate:meteo --siteaccess=project_meteo --env=prod 2>&1
*/5 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console project:flux:nrn --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project-xml-nrn.log 2>&1
13 22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:footer-encrypt --siteaccess=project --env=prod 2>&1
17 22 * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:header-encrypt --siteaccess=project --env=prod 2>&1
*/5 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console project:flux:innovation --siteaccess=project --env=prod 2>&1

0/10 * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_admin generate_podcasts 2>&1 | logify >>$LOGDIR/project-generate_podcasts.log 2>&1
0 1 * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php check_podcasts_adserving 2>&1 | logify >>$LOGDIR/project-check_podcasts_adserving.log 2>&1
0 * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS -d 'memory_limit=256M' runcronjobs.php generate_podcasts_adserving 2>&1 | logify >>$LOGDIR/project-generate_podcasts_adserving.log 2>&1
*/10 * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php unstack_podcasts_tracking_pending 2>&1 | logify >>$LOGDIR/project-unstack_podcasts_tracking_pending.log 2>&1

* * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project_admin videoUpload --service=Arkena 2>&1
*/2 * * * *    www-data    cd $DOCROOT/ezpublish_legacy && umask 002 && $PHP $PHPARGS runcronjobs.php -s project generate_static_flux_ipad 2>&1 | logify >>$LOGDIR/project-generate_static_flux_ipad.log 2>&1

*/5 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console project:export:mediabong --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project-export-mediabong.log 2>&1

*/10 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console lelab:generate:tagsHeaderFooter --encrypted --siteaccess=project_lelab --env=prod 2>&1 | logify >>$LOGDIR/lelab-generate_html_part_encrypted.log 2>&1
*/10 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:blocks --siteaccess=project header --env=prod 2>&1 | logify >>$LOGDIR/project_sponso_tag.log 2>&1
27,57 * * * *    www-data    cd $DOCROOT && umask 002 && $PHP $PHPARGS ezpublish/console front:generate:frontpage-submenus --siteaccess=project --env=prod 2>&1 | logify >>$LOGDIR/project_frontpage_submenu.log 2>&1

